# pacbuild/testsuite/testAbs/base/dcron/crontab
#
# /etc/crontab
#
# DO NOT EDIT THIS FILE MANUALLY!! USE crontab -e INSTEAD.
#

# <minute> <hour> <day> <month> <dow> <command>

00 20 * * *       /etc/cron/cron.daily
05 20 * * mon     /etc/cron/cron.weekly
10 20 1 * *       /etc/cron/cron.monthly

# End of file
