# pacbuild/testsuite/testAbs/base/dcron/cron.monthly
#!/bin/sh
#
# /etc/cron/cron.monthly
#

. /etc/cron/rotatelog

rotatelog wtmp

# End of file

